# ops/launchd/cataloged/com.adamnoonan.quant-work.lvr-pipeline-once.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.quant-work.lvr-pipeline-once</string>

    <!--
      Python entrypoint, not bash: macOS launchd's TCC sandbox blocks
      /bin/bash from reading scripts under ~/Documents (verified
      2026-04-26 — exit 126 "Operation not permitted"). Python under the
      same path works fine, mirroring the GeckoTerminal fetcher pattern.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Documents/quant-work/.venv/bin/python</string>
        <string>-m</string>
        <string>lvr.run_pipeline</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Documents/quant-work</string>

    <!--
      Fire immediately on load, then every 30 min. Avoids the local-time
      vs UTC ambiguity of StartCalendarInterval. Post-LVR-v0.7-cutover
      (April 2026) the runner is a pure consumer of scryer parquet: it
      pre-flights the required scryer datasets, runs lvr_calc ->
      coverage -> plot_calibration, writes a sentinel, and self-unloads
      via `launchctl bootout`. If scryer parquet is missing the runner
      logs the populate command and exits non-zero so a later fire (or
      a manual run after the populate) can retry.
    -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/quant-work/lvr-pipeline-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/quant-work/lvr-pipeline-launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
