# ops/launchd/com.adamnoonan.scryer.runner-backed-nav-strikes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.scryer.runner-backed-nav-strikes</string>

    <!--
        M3.6 Phase B per-manifest plist. Fires `scryer-runner tick
        --only backed-nav-strikes` every 60s. Per-manifest plist
        (vs. the multi-manifest `runner-tick`) so launchd's natural
        skip-if-running behavior provides concurrency control: if a
        previous fire is still running, the next 60s tick is
        skipped — preferable to spawning a duplicate.

        Legacy `com.adamnoonan.scryer.backed-nav-strikes.plist`
        remains loaded for parallel-soak; operator unloads it after
        verifying the runner-driven pipeline writes
        `xstocks_api_v2:runner`-attributed rows on cadence.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin/scryer-runner</string>
        <string>tick</string>
        <string>--only</string>
        <string>backed-nav-strikes</string>
        <string>--manifests</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/manifests</string>
        <string>--dataset</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/dataset</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Library/Application Support/scryer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-backed-nav-strikes.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-backed-nav-strikes.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <!-- See com.adamnoonan.scryer.runner-tick.plist for rationale. -->
        <key>PATH</key>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOSTNAME</key>
        <string>samachi-mac</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
