# ops/launchd/com.adamnoonan.scryer.runner-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.scryer.runner-tick</string>

    <!--
        M3.4 parallel-soak runner. Ticks every 60s; per-manifest
        sensors (e.g. `interval(3600s)` for kraken-trades) decide
        when to actually fire. Multi-manifest by design: future
        soak waves (M3.5+) onboard manifests to ops/sources/ rather
        than adding new launchd plists.
        Loads every *.toml in the staged manifests directory argument
        below (incl. yahoo-single-stock-iv wishlist 52: daily ATM IV snapshot,
        sensor interval 86400s — not skipped; sourced from ops/sources via
        scryer deploy).
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin/scryer-runner</string>
        <string>tick</string>
        <!--
            M3.6: skip manifests with their own dedicated per-manifest
            plist. Both plists share the runner's state file, so the
            sensor evaluator alone almost always avoids double-fires;
            but the read-state / save-state window is non-zero, and
            high-cadence manifests would race here. Listing each
            phase-B manifest by name keeps the multi-manifest tick
            handling only the low-cadence (phase A) manifests.
        -->
        <string>--skip</string>
        <string>pyth-tape</string>
        <!--
            jito-bundle-tape (wishlist 51a) — 60s sensor, ~80-150s
            wall-clock per fire. Per-manifest plist
            `com.adamnoonan.scryer.runner-jito-bundle-tape` owns it;
            skipping here so the multi-manifest tick isn't blocked by
            the multi-MB getBlock walks.
        -->
        <string>--skip</string>
        <string>jito-bundle-tape</string>
        <!-- 51e dex-xstock-swaps gets its own dedicated plist for the
             same bursty-volume reason. -->
        <string>--skip</string>
        <string>dex-xstock-swaps</string>
        <!--
            M3.6 Phase B (2026-05-03): backed-nav-strikes,
            cex-stock-perp-tape, kamino-scope-tape, v5-tape each got
            their own dedicated `runner-<id>.plist`. Skip them here
            so the multi-manifest tick doesn't race the per-manifest
            tick on shared per-manifest state.
        -->
        <string>--skip</string>
        <string>backed-nav-strikes</string>
        <string>--skip</string>
        <string>cex-stock-perp-tape</string>
        <string>--skip</string>
        <string>kamino-scope-tape</string>
        <string>--skip</string>
        <string>v5-tape</string>
        <!--
            Phase-A (remain on this plist): manifests not listed above
            are evaluated each tick — including yahoo-single-stock-iv
            (wishlist 52, daily IV via interval 86400s). Ensure scryer
            deploy synced ops/sources/yahoo-single-stock-iv.toml into
            ~/Library/Application Support/scryer/manifests so forward
            accumulation runs.
        -->
        <string>--manifests</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/manifests</string>
        <string>--dataset</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/dataset</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Library/Application Support/scryer</string>

    <key>RunAtLoad</key>
    <true/>

    <!--
        60-second tick. Sensor evaluation is cheap (no I/O for
        interval/daily). For interval(3600s) manifests this means at
        most one fire per hour with ~60s scheduling jitter. Daily and
        long-interval sensors (analytics daily(), yahoo-single-stock-iv
        86400s, etc.) are evaluated each tick too; scheduling jitter ≤60s.
    -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-tick.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-tick.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <!--
            scry is staged in Application Support; launchd's default
            PATH does not include user-staged dirs, so the runner's
            spawn of `scry` (by name, per the manifest `command` lock)
            would otherwise fail. Prepend the staging dir.
        -->
        <key>PATH</key>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <!--
            Stamps the `runner_host` column on every workflow_run row.
            macOS launchd does not export $HOSTNAME by default, so set
            it explicitly here rather than at the runner's hostname
            fallback (which would otherwise resolve to `unknown-host`).
        -->
        <key>HOSTNAME</key>
        <string>samachi-mac</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
