# ops/launchd/com.adamnoonan.scryer.runner-v5-tape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.scryer.runner-v5-tape</string>

    <!--
        M3.6 Phase B per-manifest plist. Fires `scryer-runner tick
        --only v5-tape` every 60s. Stage-2 routing via
        `--use-get-transaction` (set inside the manifest's
        [fetch].args) so the joined Chainlink + Jupiter walk goes
        through the local proxy across the healthy non-Helius
        providers. Sequential per-sig — slower than Helius
        parseTransactions but quota-resilient. Pair with the
        manifest's `--lookback-secs 120` so the fire fits inside
        the 60s cadence.

        Legacy `com.adamnoonan.scryer.v5-tape.plist` remains loaded
        for parallel-soak.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin/scryer-runner</string>
        <string>tick</string>
        <string>--only</string>
        <string>v5-tape</string>
        <string>--manifests</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/manifests</string>
        <string>--dataset</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/dataset</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Library/Application Support/scryer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-v5-tape.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-v5-tape.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>PATH</key>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOSTNAME</key>
        <string>samachi-mac</string>
        <!--
            v5-tape's CLI reads HELIUS_API_KEY from the env even when
            --use-get-transaction is set (the URL is built but unused).
            Set this in your ~/.zshenv or via launchctl setenv.
        -->
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
