# ops/launchd/com.adamnoonan.scryer.runner-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.scryer.runner-tick</string>

    <!--
        M3.4 parallel-soak runner. Ticks every 60s; per-manifest
        sensors (e.g. `interval(3600s)` for kraken-trades) decide
        when to actually fire. Multi-manifest by design: future
        soak waves (M3.5+) onboard manifests to ops/sources/ rather
        than adding new launchd plists.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin/scryer-runner</string>
        <string>tick</string>
        <string>--manifests</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/manifests</string>
        <string>--dataset</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/dataset</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Library/Application Support/scryer</string>

    <key>RunAtLoad</key>
    <true/>

    <!--
        60-second tick. Sensor evaluation is cheap (no I/O for
        interval/daily). For interval(3600s) manifests this means at
        most one fire per hour with ~60s scheduling jitter. For daily
        sensors (none yet) it bounds the worst-case lateness to 60s.
    -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-tick.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-tick.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <!--
            scry is staged in Application Support; launchd's default
            PATH does not include user-staged dirs, so the runner's
            spawn of `scry` (by name, per the manifest `command` lock)
            would otherwise fail. Prepend the staging dir.
        -->
        <key>PATH</key>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <!--
            Stamps the `runner_host` column on every workflow_run row.
            macOS launchd does not export $HOSTNAME by default, so set
            it explicitly here rather than at the runner's hostname
            fallback (which would otherwise resolve to `unknown-host`).
        -->
        <key>HOSTNAME</key>
        <string>samachi-mac</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
