# ops/launchd/com.adamnoonan.scryer.runner-pyth-lazer-tape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adamnoonan.scryer.runner-pyth-lazer-tape</string>

    <!--
        Per-manifest plist. Polls `scryer-runner tick --only
        pyth-lazer-tape` every 10s; the manifest sensor `interval(60s)`
        in `ops/sources/pyth-lazer-tape.toml` enforces the actual
        60s fire cadence — most launchd ticks therefore evaluate the
        sensor and `Hold` cheaply (sub-100ms), with one fire per
        60-second window.

        Why poll-and-hold instead of launchd's `StartInterval=60`:
        macOS launchd's `StartInterval=N` semantic empirically waits
        `(previous_job_exit + N)`, not `(previous_job_start + N)`. With
        per-fire wall-clock ~41s (30s WS drain + ~11s parquet merge-
        dedup across the 26-feed panel), `StartInterval=60` produced
        a 101s effective cadence instead of 60s. `StartInterval=10`
        with skip-if-running and the in-binary `interval(60s)` sensor
        recovers the 60s cadence: each fire blocks the next 3-4
        launchd ticks (skip-if-running), then the first poll after
        the job exits and 60s-elapsed checks both pass and fires.

        No legacy plist sibling — Pyth Lazer is a new ingest path
        landing 2026-05-10, so this is the only deploy target.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin/scryer-runner</string>
        <string>tick</string>
        <string>--only</string>
        <string>pyth-lazer-tape</string>
        <string>--manifests</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/manifests</string>
        <string>--dataset</string>
        <string>/Users/adamnoonan/Library/Application Support/scryer/dataset</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/adamnoonan/Library/Application Support/scryer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-pyth-lazer-tape.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamnoonan/Library/Logs/scryer/runner-pyth-lazer-tape.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>PATH</key>
        <string>/Users/adamnoonan/Library/Application Support/scryer/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOSTNAME</key>
        <string>samachi-mac</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
